# templates/geth.service.erb
[Unit]
Description=Go Ethereum Client
After=network.target
Wants=network.target

[Service]
User=goeth
Group=goeth
Type=simple
Restart=always
RestartSec=5
ExecStart= <%= executable %> <%= flag %> --http --datadir <%= datadir %> --cache 2048 --maxpeers 30

[Install]
WantedBy=default.target
