# templates/geth.service.erb
[Unit]
Description=Go Ethereum Client
After=network.target
Wants=network.target

[Service]
User=goeth
Group=goeth
Type=simple
Restart=always
RestartSec=5
ExecStart= <%= executable %> <%= flag %> --http --datadir <%= datadir %> --port <%= port %> --cache 8192 --maxpeers 50

[Install]
WantedBy=default.target
